# moeka.service
[Unit]
Description=Moeka (nanobot) agent gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/projects/moeka
# 'run' starts the gateway in the foreground so systemd owns the PID.
# Do NOT use 'start' here — that backgrounds the process and breaks restart.
ExecStart=%h/projects/moeka/moeka.sh run
ExecStop=/bin/kill -TERM $MAINPID
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
# Give the process up to 30s to shut down gracefully before SIGKILL
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
# moeka.sh sources ./.env and ./keys.env itself; these lines are belt-and-braces
# for cases where the script is invoked without them (e.g. `systemctl edit`).
EnvironmentFile=-%h/projects/moeka/.env
EnvironmentFile=-%h/projects/moeka/keys.env

[Install]
WantedBy=default.target
