# scripts/moeka.service
[Unit]
Description=Moeka (nanobot) agent gateway
After=network-online.target
Wants=network-online.target
# Never give up — allow unlimited restart attempts.
StartLimitBurst=0

[Service]
Type=simple
WorkingDirectory=%h/projects/moeka
# 'run' starts the gateway in the foreground so systemd owns the PID.
# Do NOT use 'start' here — that backgrounds the process and breaks restart.
ExecStart=%h/projects/moeka/bin/moeka.sh run
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
# Exponential backoff: 5s → ~10s → ~20s → ~40s → ~80s → 5min (cap).
# Requires systemd 254+.
RestartSec=5
RestartSteps=6
RestartMaxDelaySec=5min
# Give the process up to 30s to shut down gracefully before SIGKILL
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
# moeka.sh sources ./.env and ./keys.env itself; these lines are belt-and-braces
# for cases where the script is invoked without them (e.g. `systemctl edit`).
EnvironmentFile=-%h/projects/moeka/.env
EnvironmentFile=-%h/projects/moeka/keys.env

[Install]
WantedBy=default.target
